# root/etc/services.d/lms/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid dockeruser /usr/sbin/squeezeboxserver --prefsdir $SQUEEZE_VOL/prefs --logdir $SQUEEZE_VOL/logs --cachedir $SQUEEZE_VOL/cache
